# data/cron/crontabs/example.cron
* * * * * /timestamp.sh >> /var/log/cron/cron.log